# scripts/einig-oder-aus.service
[Unit]
Description=Einig oder Aus! — Die Demokratiemaschine
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__REPO_DIR__
ExecStart=/usr/bin/env uv run python app.py
Restart=always
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
