# scripts/einig-hotspot.service
[Unit]
Description=Konsensomat — Fallback WiFi Hotspot
After=NetworkManager-wait-online.service
Wants=NetworkManager-wait-online.service

[Service]
Type=oneshot
RemainAfterExit=yes
User=__USER__
ExecStart=__REPO_DIR__/scripts/hotspot.sh

[Install]
WantedBy=multi-user.target
